# target/linux/rockchip/armv8/base-files/etc/uci-defaults/13_pon_config
#!/bin/sh
board=$(cat /tmp/sysinfo/board_name)
case "$board" in
aio,aio-3399b)
uci -q set network.lan.ipaddr="192.168.9.1"
uci -q set network.pon=interface
uci -q set network.pon.proto="static"
uci -q set network.pon.device="eth0.2"
uci -q set network.pon.delegate="0"
uci -q set network.pon.ipaddr="192.168.1.254"
uci -q set network.pon.netmask="255.255.255.0"
uci -q set network.pon.metric="100"
uci -q set network.pon_route=route
uci -q set network.pon_route.interface='pon'
uci -q set network.pon_route.target='192.168.1.10'
uci -q set network.pon_route.netmask='255.255.255.0'
uci -q set network.pon_route.gateway='192.168.1.254'
uci -q set network.pon_route.metric='100'
uci -q commit network
chmod +x /bin/aio_pon
sed -i '2a\aio_pon' /etc/rc.local
;;
*)
exit 0
;;
esac
exit 0